# backup.service
[Unit]
After=local-fs.target

[Service]
ExecStart=/home/philip/scripts/backup-sync.sh

[Install]
WantedBy=multi-user.target
